# InfiniteBotServerRunner.service
[Unit]
Description=The Infinite Experiment bot

[Service]
User=root
WorkingDirectory=/root/projects/infinite-experiment-bot
ExecStart=/usr/bin/sh /root/projects/infinite-experiment-bot/runner.sh
Restart=always

[Install]
WantedBy=multi-user.target

